# src/mandroid.service
[Unit]
Description=The Mandroid Program

[Service]
Type=simple
ExecStart=/usr/bin/mandroid

[Install]
WantedBy=multi-user.target
